# root/etc/network/if-up.d/pistar-motdgen
#!/bin/bash
mount -o remount,rw /
/etc/init.d/procps restart
/usr/local/sbin/pistar-motdgen
