# root/var/lib/dpkg/info/minissdpd.list
/.
/etc
/etc/init.d
/etc/init.d/minissdpd
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/minissdpd.service
/usr
/usr/lib
/usr/lib/minissdpd
/usr/lib/minissdpd/minissdpd-systemd-wrapper
/usr/sbin
/usr/sbin/minissdpd
/usr/share
/usr/share/doc
/usr/share/doc/minissdpd
/usr/share/doc/minissdpd/changelog.Debian.gz
/usr/share/doc/minissdpd/changelog.gz
/usr/share/doc/minissdpd/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/minissdpd
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/minissdpd.1.gz
/usr/share/minissdpd
/usr/share/minissdpd/minissdpd.default
/etc/default/minissdpd
